# scripts/service/com.aigenius.octoally.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.aigenius.octoally</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>-c</string>
    <string>(/usr/local/bin/node -e "require('better-sqlite3')" 2&gt;/dev/null || /usr/local/bin/npm rebuild better-sqlite3) &amp;&amp; exec /usr/local/bin/node __OCTOALLY_DIR__/server/dist/index.js</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__OCTOALLY_DIR__/server</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>NODE_ENV</key>
    <string>production</string>
    <key>PORT</key>
    <string>42010</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>StandardOutPath</key>
  <string>__OCTOALLY_DIR__/logs/octoally.log</string>
  <key>StandardErrorPath</key>
  <string>__OCTOALLY_DIR__/logs/octoally-error.log</string>

  <key>UserName</key>
  <string>__USER__</string>
</dict>
</plist>
